# app-emulation/vmware-service/files/vmware.service
# It's not recommended to modify this file in-place, because it will be
# overwritten during package upgrades. If you want to customize, the
# best way is to create file
# "/etc/systemd/system/vmware.service.d/*.conf"
# containing your changes

[Unit]
Description=VMware Service
After=network.target

[Service]
ExecStart=/bin/bash -c "/opt/vmware/bin/vmware-usbarbitrator  && modprobe -vq vmci && modprobe -av vmmon vsock vmblock vmnet && /opt/vmware/bin/vmware-networks --start  | /usr/bin/dialog"  > /dev/console 2>&1
#ExecStart=/bin/bash -c "/bin/bash -c \"/opt/vmware/bin/vmware-usbarbitrator  && modprobe -vq vmci && modprobe -av vmmon vsock vmblock vmnet && /opt/vmware/bin/vmware-networks --start\" | /usr/bin/dialog --clear --backtitle \"VMware Start\"  --title \"VMware Start\" --progressbox 20 70 > /dev/console 2>&1"

ExecStop=/bin/bash -c "killall --wait /opt/vmware/bin/vmware-usbarbitrator && /opt/vmware/bin/vmware-networks --stop"
ExecReload=/bin/bash -c "killall --wait /opt/vmware/bin/vmware-usbarbitrator && /opt/vmware/bin/vmware-networks --stop"

[Install]
WantedBy=multi-user.target