# app-emulation/vmware-service/files/vmware.service
# It's not recommended to modify this file in-place, because it will be
# overwritten during package upgrades. If you want to customize, the
# best way is to create file
# "/etc/systemd/system/vmware.service.d/*.conf"
# containing your changes

[Unit]
Description=VMware Service
After=network.target

[Service]
ExecStart=/opt/vmware/bin/vmware-usbarbitrator  && modprobe -vq vmci && modprobe -av vmmon vsock vmblock vmnet && /opt/vmware/bin/vmware-networks --start
ExecStop=killall --wait /opt/vmware/bin/vmware-usbarbitrator && /opt/vmware/bin/vmware-networks --stop
ExecReload=killall --wait /opt/vmware/bin/vmware-usbarbitrator && /opt/vmware/bin/vmware-networks --stop

[Install]
WantedBy=multi-user.target