# net-misc/cow/files/cow.service
[Unit]
Description=Cow Proxy Service
Wants=network-online.target
After=network-online.target

[Service]
ExecStart=/usr/bin/cow -rc=/etc/cow/rc

[Install]
WantedBy=default.target
