# app-misc/dradis/files/dradis.service
# It's not recommended to modify this file in-place, because it will be
# overwritten during package upgrades. If you want to customize, the
# best way is to create file
# "/etc/systemd/system/dradis.service.d/*.conf"
# containing your changes

# For example, if you want to change some daemon and/or unit options,
# create a file named
# "/etc/systemd/system/dradis.service.d/myopts.conf"
# containing:
#       [Service]
#       Environment="BITCOIN_OPTS=-debug -logtimestamps"
#       Nice=10
# This will override the setting appearing below.

# Note that almost all daemon options could be specified in
# /etc/dradis/dradis.conf

[Unit]
Description=Dradis Daemon
After=network.target

[Service]
#User=dradis
#EnvironmentFile=-/etc/dradis/dradis.conf
#ExecStart=dradis --start ${EXTRA_PARAMS} \
#		--pidfile ${PID} \
#		--exec ${RUBY} ${SERVER} -- \
#		--daemon --environment=${ENVIRONMENT} \
#		--binding=${ADDRESS} --port=${PORT}
#ExecStop=dradis --stop ${EXTRA_PARAMS} \
#		--pidfile ${PID}
ExecStart="/bin/sh -c 'cd /usr/share/dradis && ruby bin/setup && bundle exec rails server'"
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target