# net-analyzer/metasploit-service/files/msfrpcd.service
# It's not recommended to modify this file in-place, because it will be
# overwritten during package upgrades. If you want to customize, the
# best way is to create file
# "/etc/systemd/system/msfrpcd.service.d/*.conf"
# containing your changes

# Note that almost all daemon options could be specified in
# /etc/msf/msf.conf

[Unit]
Description=MSF RPC Daemon
After=network.target

[Service]
EnvironmentFile=-/etc/msf/msfrpcd.conf
ExecStart=/bin/bash -c "cd /usr/lib/metasploit/ && /usr/lib/metasploit/msfrpcd ${MSF_OPTS}"
ExecStop=/bin/kill -HUP $MAINPID
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target