# net-proxy/shadowsocks-service/files/sslocal.service
# It's not recommended to modify this file in-place, because it will be
# overwritten during package upgrades. If you want to customize, the
# best way is to create file
# "/etc/systemd/system/sslocal.service.d/*.conf"
# containing your changes

[Unit]
Description=Shadowsocks Client Daemon
After=network.target

[Service]
ExecStart=/bin/bash -c "sslocal -c /etc/shadowsocks.json -d start"
ExecStop=/bin/bash -c "sslocal -c /etc/shadowsocks.json -d stop"
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target