# www-apps/spiderfoot/files/spiderfoot.service
# It's not recommended to modify this file in-place, because it will be
# overwritten during package upgrades. If you want to customize, the
# best way is to create file
# "/etc/systemd/system/spiderfoot.service.d/*.conf"
# containing your changes

# Note that almost all daemon options could be specified in
# /etc/spiderfoot/spiderfoot.conf

[Unit]
Description=SpiderFoot Daemon
After=network.target

[Service]
EnvironmentFile=-/etc/spiderfoot/spiderfoot.conf
ExecStop=/bin/kill -HUP $MAINPID
ExecStart=/bin/bash -c "python /usr/share/spiderfoot/sf.py ${ADDRESS}:${PORT} ${EXTRA_PARAMS}"
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target