# www-apps/spiderfoot/files/spiderfoot.initd
#!/sbin/runscript
# Copyright 1999-2010 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

depend() {
	need net
}

start() {
	ebegin "Starting spiderfoot"
	python /usr/share/spiderfoot/sf.py ${ADDRESS}:${PORT} ${EXTRA_PARAMS}
	eend $?
}

stop() {
	ebegin "Stopping spiderfoot"
	killall spiderfoot
	eend $?
}
